# sbd.service
[Unit]
Description=Start sbd
After=syslog.target network.target
[Service]
User=root
Environment=PATH='/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin'
Type=forking
ExecStart=/usr/sbin/sbd start >> /var/log/sbd.log 2>&1
ExecStop=/usr/sbin/sbd stop >> /var/log/sbd.log 2>&1
ExecReload=/usr/sbin/sbd restart >> /var/log/sbd.log 2>&1 
StandardOutput=tty 
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
